# server-monitor-bot.service
[Unit]
Description=Telegram Server Monitor Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/USER/server-monitor-bot
ExecStart=/home/USER/server-monitor-bot/venv/bin/python bot.py
Restart=always
RestartSec=5
Environment="BOT_TOKEN=your-token"
Environment="OWNER_ID=123456789"
Environment="OWNER_USERNAME=tg-username-without-@"

# Security hardening (Optional)
User=USER
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
